# usr/lib/systemd/user/laptop_battery_status.timer
[Unit]
Description=Check battery status script every 10 minutes

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
