# usr/lib/systemd/user/laptop_battery_status.service
[Unit]
Description=Laptop battery status

[Service]
ExecStart=/home/ccw/.local/bin/laptop_battery_status.sh
